# run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
After=blockdev@dev-disk-by\x2duuid-386769a3\x2dcfa5\x2d47c8\x2d8797\x2dd5ec58c9cb6c.target

[Mount]
Where=/
What=/dev/disk/by-uuid/386769a3-cfa5-47c8-8797-d5ec58c9cb6c
Type=ext4
