# run/systemd/generator/mnt-jupiter_1.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
After=blockdev@dev-disk-by\x2did-scsi\x2d0DO_Volume_jupiter\x2d1.target

[Mount]
Where=/mnt/jupiter_1
What=/dev/disk/by-id/scsi-0DO_Volume_jupiter-1
Type=ext4
Options=defaults,nofail,discard
